# android/gmbot/gmbot/etc/init/uwsgi.conf
description     "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]
respawn

env BINPATH=/usr/local/sbin/uwsgi
exec $BINPATH --uid nginx --ini /etc/uwsgi/uwsgi.ini --pidfile /var/run/uwsgi.pid

